# app/kitchen/site-cookbooks/notify-email-on-boot/templates/default/notify-email-on-boot.conf.erb
start on runlevel [!06]

task

script

PUBLIC_IP=`curl --connect-timeout 5 -s icanhazip.com`

MAILADD=<%= @mail_to %>

SERVER=`hostname --fqdn`

MESSAGE=$(cat <<EOF
$SERVER has booted up.

Public IP: $PUBLIC_IP

If this is news to you, please investigate.

`date -u`
EOF
)

exec echo "$MESSAGE" | mail -s "Boot of $SERVER" $MAILADD

end script
