# app/kitchen/site-cookbooks/docker-containers/templates/default/docker-rpc.conf.erb
description "RSM-DATA RPC Docker Container"
author "systemsadmin@ridesharemarket.com"
start on filesystem and started docker
stop on runlevel [!2345]
respawn
script
  /usr/bin/docker start rpc
end script
